# provisioners/webapp.service
[Unit]
Description=webapp API
ConditionPathExists=/opt/webapp/.env
After=network.target

[Service]
Type=simple
User=csye6225
Group=csye6225
WorkingDirectory=/opt/webapp
ExecStart=/opt/webapp/venv/bin/python3 /opt/webapp/app.py
Restart=always
RestartSec=3
StandardOutput=append:/var/log/webapp.log
StandardError=append:/var/log/webapp.log

[Install]
WantedBy=multi-user.target
